# deploy/hifleet-agent.service
[Unit]
Description=HiFleet Agent HTTP Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ecs-user
Group=ecs-user
WorkingDirectory=/home/ecs-user/coze_ai
ExecStart=/bin/bash /home/ecs-user/coze_ai/scripts/start_system_service.sh
Restart=always
RestartSec=3
LimitNOFILE=65535
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
